# debian/cron.weekly
#!/bin/sh -e
cd /home/judge
find . -name "[MTP]*" -atime +90 -exec rm {} \;
./flock ./dip.log > /dev/null
./newlogs
./flock > /dev/null
